# packer/include/zk-exhibitor.conf
description "zk-exhibitor container"

start on started docker
stop on stopping docker

script
  /usr/bin/docker start -a zk-exhibitor
end script

respawn
respawn limit 10 10
kill timeout 10